# devops-swarm-project/app/web/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm8.2 -F

[program:nginx]
command=nginx -g "daemon off;"
